# ansible/roles/backup/templates/restic-backup.timer.j2
[Unit]
Description=Daily Restic Backup Timer for Mediacenter

[Timer]
OnCalendar={{ restic_timer_on_calendar }}
RandomizedDelaySec={{ restic_timer_randomized_delay }}
Persistent=true

[Install]
WantedBy=timers.target
